# docker/vulnbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
autostart=true
autorestart=true
environment=APACHE_RUN_USER=www-data,APACHE_RUN_GROUP=www-data,APACHE_RUN_DIR=/var/run/apache2,APACHE_PID_FILE=/var/run/apache2/apache2.pid,APACHE_LOCK_DIR=/var/lock/apache2,APACHE_LOG_DIR=/var/log/apache2

[program:mysql]
command=/usr/bin/mysqld_safe --port=3306
autostart=true
autorestart=true

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true

[program:snmpd]
command=/usr/sbin/snmpd -f -Lo -u root -c /etc/snmp/snmpd.conf udp:161
autostart=true
autorestart=true

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autostart=true
autorestart=true

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autostart=true
autorestart=true

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true

[program:vulnapi]
command=python3 /opt/vulnapi/app.py
autostart=true
autorestart=true
user=root
environment=FLASK_PORT=9090

# Redis — no auth, bound to all interfaces
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
user=redis
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log

# BIND9 DNS — open resolver with vulnerable zones
[program:named]
command=/usr/sbin/named -f -u bind
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/named.log
stderr_logfile=/var/log/supervisor/named.log

# OpenLDAP — anonymous bind, weak passwords
[program:slapd]
command=/usr/sbin/slapd -d 0 -h "ldap:/// ldapi:///"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/slapd.log
stderr_logfile=/var/log/supervisor/slapd.log

# Postfix SMTP — open relay
[program:postfix]
command=/bin/bash -c "postfix start-fg"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postfix.log
stderr_logfile=/var/log/supervisor/postfix.log
