# docker/vulnbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
autostart=true
autorestart=true
environment=APACHE_RUN_USER=www-data,APACHE_RUN_GROUP=www-data

[program:mysql]
command=/usr/bin/mysqld_safe --port=3306
autostart=true
autorestart=true

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true

[program:snmpd]
command=/usr/sbin/snmpd -f -Lo -u root -c /etc/snmp/snmpd.conf udp:161
autostart=true
autorestart=true

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autostart=true
autorestart=true

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autostart=true
autorestart=true

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true

[program:vulnapi]
command=python3 /opt/vulnapi/app.py
autostart=true
autorestart=true
user=root
environment=FLASK_PORT=9090

# Redis — no auth, bound to all interfaces (vuln: #1)
[program:redis]
command=/usr/sbin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
user=redis
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log
